# ansible/roles/wakeonlan/templates/wol.service.j2
[Unit]
Description=Configure Wake-up on LAN

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s {{ eth_nic }} wol g

[Install]
WantedBy=basic.target
